# rootfs/etc/systemd/system/ckan.service.d/20-wait-solr.conf
[Service]
ExecStartPre=/bin/bash -lc 'for i in {1..60}; do curl -fsS http://127.0.0.1:8983/solr/ckan/admin/ping >/dev/null && exit 0; sleep 1; done; echo "Solr (core ckan) não respondeu ao ping em 60s" >&2; exit 1'
