# rootfs/etc/systemd/system/ckan.service.d/99-bind-local.conf
[Service]
ExecStart=
ExecStart=/usr/lib/ckan/venv/bin/ckan run --host 127.0.0.1 --port 5000
