# rootfs/etc/systemd/system/ckan.service.d/solr.conf
[Unit]
Wants=solr.service
After=solr.service
